# data/bbswitchd.service.in
[Unit]
Description=GPU power state management daemon

[Service]
Type=simple
ExecStart=@SBINDIR@/bbswitchd --log=syslog --fd=3
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
